# wipe-monthly-full.service
[Unit]
Description=Monthly full wipe

[Service]
Type=oneshot
ExecStart=/home/rust/wipe-monthly-full.sh
