# wipe-weekly-map-only.timer
[Unit]
Description=Weekly map wipe

[Timer]
OnCalendar=Fri *-*-* 14:00:00
Persistent=true

[Install]
WantedBy=timers.target
