# wipe-weekly-map-only.service
[Unit]
Description=Weekly map wipe

[Service]
Type=oneshot
ExecStart=/home/rust/wipe-weekly-map-only.sh
